# srv/openbsd/etc/rc.d/urloli
#!/bin/ksh
#
# $OpenBSD: urloli.rc,v 1.4 2018/01/11 19:27:11 rpe Exp $

name="urloli"
daemon="/usr/local/bin/${name}"

. /etc/rc.d/rc.subr

rc_cmd $1
